# deploy/syslog-ng/build/syslog-ng.ini
[program:syslog-ng]
command=/usr/sbin/syslog-ng -Fvde -f /etc/syslog-ng/syslog-ng.conf